# undervolt.service
[Unit]
Description=undervolt

[Service]
Type=oneshot
ExecStart=/usr/local/bin/undervolt -v --core -105 --cache -105 --gpu -70
